# systemd/watch-ip.timer
[Unit]
Description=Run the check-ip script every 2 minutes.

[Timer]
OnUnitActiveSec=2m

[Install]
WantedBy=timers.target
